# systemd/timer/hello.timer
[Unit]
Description=Run hello every 5 min

[Timer]
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
